# working/profile.d/debuginfod.sh
/usr/share/libdebuginfod-common/debuginfod.sh