# working/profile.d/debuginfod.csh
/usr/share/libdebuginfod-common/debuginfod.csh