# working/profile.d/tmout.sh

# Set TMOUT to 900 per security requirements
TMOUT=900
readonly TMOUT
export TMOUT
